# workshop-forge.service
[Unit]
Description=Workshop Forge Voice Assistant
Documentation=https://github.com/anthropics/workshop-forge
After=network-online.target sound.target
Wants=network-online.target

[Service]
Type=simple
User=tyler
Group=tyler
WorkingDirectory=/home/tyler/projects/workshop-assistant
Environment=PATH=/home/tyler/projects/workshop-assistant/venv/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=/home/tyler/projects/workshop-assistant/venv/bin/python src/main.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Audio permissions
SupplementaryGroups=audio

# Graceful shutdown
TimeoutStopSec=15
KillMode=mixed

[Install]
WantedBy=multi-user.target